# ops/macos/com.archivox.web.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.archivox.web</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/zsh</string>
      <string>__ARCHIVOX_ROOT__/scripts/run_archivox.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__ARCHIVOX_ROOT__</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
      <key>UV_BIN</key>
      <string>/opt/homebrew/bin/uv</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/archivox-web.out.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/archivox-web.err.log</string>
  </dict>
</plist>
